# systemback-ng/systemd/systemback-ng.service
[Unit]
Description=Systemback-ng scheduled backup

[Service]
Type=oneshot
ExecStart=/usr/bin/systemback-ng backup --label auto
